# deploy/flyio/s6/llm/run
#!/bin/execlineb -P
# s6 service script for LLM inference server
# Port 8080

s6-envdir /etc/s6-overlay/env
fdmove -c 2 1
tri_inference
