# fpga/tools/trinity-ralph.service
# ═══════════════════════════════════════════════════════════════════════════════
# TRINITY RALPH ORCHESTRATOR v2.0 — 24/7 Autonomous AI Service
# ═══════════════════════════════════════════════════════════════════════════════
#
# Runs 24/7 autonomous self-improvement cycles:
# - Every 10 min: Task processing (fix_plan.md)
# - Every 30 min: Research ingestion (arXiv papers)
# - Every 2 hours: Full evolution cycle (improve → test → flash → publish)
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp trinity-ralph.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable trinity-ralph.service
#   systemctl --user start trinity-ralph.service
#
# Usage:
#   systemctl --user status trinity-ralph
#   systemctl --user stop trinity-ralph
#   journalctl --user -u trinity-ralph -f
#
# ═══════════════════════════════════════════════════════════════════════════════

[Unit]
Description=Trinity Ralph Orchestrator — 24/7 Autonomous Self-Evolution AI
Documentation=https://github.com/gHashTag/trinity
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/Users/playra/trinity-w1/fpga/tools/ralph-cron.sh start
ExecStop=/Users/playra/trinity-w1/fpga/tools/ralph-cron.sh stop
Restart=always
RestartSec=60

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin"
Environment="RALPH_REPORT_ENABLED=true"
Environment="ANTHROPIC_API_KEY="  # Set via systemd credential
WorkingDirectory=/Users/playra/trinity-w1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trinity-ralph

# Resource limits (increased for 24/7 operation)
MemoryMax=1G
CPUQuota=75%

# Safety timeouts
TimeoutStartSec=300
TimeoutStopSec=60

[Install]
WantedBy=default.target
