# deploy/flyio/s6/mcp/run
#!/bin/execlineb -P
# s6 service script for MCP WebSocket server
# Port 8081

s6-envdir /etc/s6-overlay/env
fdmove -c 2 1
trinity-mcp
